# bin/ss
#! /bin/bash
### BEGIN INIT INFO
# Provides:          Shadowsocks-libev
# Short-Description: Fast tunnel proxy that helps you bypass firewalls
# Description:       Start or stop the Shadowsocks-libev client
### END INIT INFO

# Author: vinew.cc

if [ -f /usr/local/bin/ss-local ]; then
    DAEMON=/usr/local/bin/ss-local
elif [ -f /usr/bin/ss-local ]; then
    DAEMON=/usr/bin/ss-local
fi
NAME=Shadowsocks-libev
CONF=$HOME/APP/SRK/config/ss_client.json
DIR=$HOME/APP/SRK
PAC_PID=$DIR/pid/pac.pid
KCP_PID=$DIR/pid/kcptun.pid
SS_PID=$DIR/pid/shadowsocks-libev.pid
RET_VAL=0

if [ ! -d $DIR ]; then
    mkdir -p $DIR
    if [ $? -ne 0 ]; then
        echo "Creating PID directory $DIR failed"
        exit 1
    fi
fi

if [ ! -f $CONF ]; then
    echo "$NAME config file $CONF not found"
    exit 1
fi

check_PAC() {
    if [ -r $PAC_PID ]; then
        read PAC < $PAC_PID
        if ps -p $PAC | grep SimpleHTTPServer >/dev/null; then
            return 0
        else
            rm -f $PAC_PID
            return 1
        fi
    else
        return 2
    fi
}

check_kcptun() {
    if [ -r $KCP_PID ]; then
        read KCP < $KCP_PID
        if ps -p $KCP | grep client >/dev/null; then
            return 0
        else
            rm -f $KCP_PID
            return 1
        fi
    else
        return 2
    fi
}

check_ss() {
    if [ -r $SS_PID ]; then
        read SS < $SS_PID
        if ps -p $SS | grep ss-local>/dev/null; then
            return 0
        else
            rm -f $SS_PID
            return 1
        fi
    else
        return 2
    fi
}

start() {
    if check_ss; then
        echo "$NAME (pid $SS) is already running..."
        return 0
    fi
    if check_PAC; then
        echo "PAC (pid $PAC) is already running..."
    else
        bash ${DIR}/pac.sh
        echo "Starting PAC success"
    fi
    if check_kcptun; then
        echo "kcptun (pid $KCP) is already running..."
    else
        bash ${DIR}/kcp_ss.sh
        echo "Starting kcptun success"
    fi
    
    $DAEMON -c $CONF -f $SS_PID
    if check_ss; then
        echo "Starting $NAME success"
    else
        echo "Starting $NAME failed"
        RET_VAL=1
    fi
}

stop() {
    if check_PAC; then
        kill -9 $PAC
        rm -f $PAC_PID
        echo "Stopping PAC success"
    else
        echo "PAC is stopped"
    fi
    if check_kcptun; then
        kill -9 $KCP
        rm -f $KCP_PID
        echo "Stopping kcptun success"
    else
        echo "kcptun is stopped"
    fi
    if check_ss; then
        kill -9 $SS
        rm -f $SS_PID
        echo "Stopping $NAME success"
    else
        echo "$NAME is stopped"
        RET_VAL=1
    fi
}

status() {
    check_PAC
    case $? in
        0)
        echo "PAC (pid $PAC) is running..."
        ;;
        1|2)
        echo "PAC is stopped"
        RET_VAL=1
        ;;
    esac
    check_kcptun
    case $? in
        0)
        echo "kcptun (pid $KCP) is running..."
        ;;
        1|2)
        echo "kcptun is stopped"
        RET_VAL=1
        ;;
    esac
    check_ss
    case $? in
        0)
        echo "$NAME (pid $SS) is running..."
        ;;
        1|2)
        echo "$NAME is stopped"
        RET_VAL=1
        ;;
    esac
}

pacupdate() {
   bash ${DIR}/gfwlist2pac.sh
}

case "$1" in
'start')
    start
    ;;   
'stop')
    stop
    ;; 
'status')
    status
    ;;
'restart')
    stop
    start
    RETVAL=$?
    ;;  
'pacupdate')
    pacupdate
    ;; 
*)
    echo "Usage: $0 { start | stop | restart | status | pacupdate}"
    RETVAL=1
esac
exit $RETVAL